# private_dot_config/systemd/user/generateWeather.service
[Unit]
Description=Generate weather data in ~/.cache/.weather_cache

[Service]
Type=oneshot
ExecStart=/home/alberth/.local/bin/Weather.py
