# private_dot_config/systemd/user/generateWeather.timer
[Unit]
Description=Generate weather data

[Timer]
OnBootSec=10m
OnUnitActiveSec=15m

[Install]
WantedBy=timers.target
